# docs/service/display-client.service.example
[Unit]
Description=E-Paper display client
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=6h
StartLimitBurst=3

[Service]
Type=notify
NotifyAccess=main
User=pi
WorkingDirectory=/home/pi/display_programme
Environment=DISPLAY_HEALTH_PATH=/run/rpi-waiting-time-display/client-health.json
RuntimeDirectory=rpi-waiting-time-display
RuntimeDirectoryMode=0755
ExecStart=/home/pi/display_env/bin/python /home/pi/display_programme/basic.py
KillSignal=SIGTERM
TimeoutStartSec=180
TimeoutStopSec=30
WatchdogSec=360
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpi-display-client
LogRateLimitIntervalSec=30s
LogRateLimitBurst=500

[Install]
WantedBy=multi-user.target

# This unit is for the split client/server architecture.  The process named in
# ExecStart must be the actual client so NotifyAccess=main cannot be satisfied
# by a shell wrapper or an unrelated helper process.
