# docs/service/webserial.service.example
[Unit]
# Version: 0.0.3 (2025-01-12) # AUTO-INCREMENT
Description=WebSerial Configuration Interface
After=network.target bluetooth.target
Wants=network.target bluetooth.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/display_programme
ExecStart=/home/pi/display_env/bin/python3 /home/pi/display_programme/webserial_server.py
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/var/log/webserial.out
StandardError=append:/var/log/webserial.err
Restart=on-failure
RestartSec=5

# Ensure proper permissions
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p /home/pi/display_programme/logs
ExecStartPre=/bin/chown -R pi:pi /home/pi/display_programme/logs

[Install]
WantedBy=multi-user.target

# Installation instructions:
# Note: This service is automatically installed by setup_display.sh
# These instructions are only needed if you want to install it manually:
#
# 1. Make sure to change "pi" to the username of the account you are using, e.g. /home/pi --> /home/your_username
#
# 2. Copy this file to /etc/systemd/system/webserial.service:
#    sudo cp webserial.service.example /etc/systemd/system/webserial.service
#
# 3. Reload systemd:
#    sudo systemctl daemon-reload
#
# 4. Enable the service to start on boot:
#    sudo systemctl enable webserial.service
#
# 5. Start the service:
#    sudo systemctl start webserial.service
#
# 6. Check status:
#    sudo systemctl status webserial.service
#
# 7. View logs:
#    journalctl -u webserial.service -f
#    tail -f /var/log/webserial.out
#    tail -f /var/log/webserial.err